# dockerfiles/webdav/s6/nginx/run
#!/bin/bash
#set -eux

/usr/sbin/nginx -g "daemon off;"
